# config/.runit/service/gdrive/run
#!/bin/sh

exec rclone mount gdrive: ~/gdrive
