# config/.runit/service/mpd/run
#!/bin/sh

sv check "$HOME/.runit/service/pipewire" >/dev/null || exit 1
exec mpd --no-daemon
